# cron.daily/motion-cleanup
#!/bin/bash

/usr/bin/find /volumes/data/motion/data -mtime 7 -delete

